# bare-metal/services/head/ray-gpu5-tunnel.service
[Unit]
Description=SSH tunnel to gpu5

[Service]
Type=simple
Restart=always
RestartSec=10
ExecStart=/usr/bin/ssh -N -R 6379:localhost:6379 gpu5

[Install]
WantedBy=default.target
